# entryfiles/_etc/supervisor/supervisord.d/dummy1.conf
[program:dummy1]
command=sh -c "date >> /tmp/dummy1_timestamp"
startsecs = 0
autorestart = false
startretries = 1